# 8.1/alpine3.14/nginx/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

exec docker-php-entrypoint php-fpm
