# 8.0/bullseye/nginx/rootfs/etc/services.d/php-fpm/run
#!/command/with-contenv sh

exec docker-php-entrypoint php-fpm
