# 8.1/alpine3.15/nginx/rootfs/etc/s6-overlay/s6-rc.d/nginx-conf-envsubst/up
with-contenv /opt/nginx/nginx-conf-envsubst.sh
