# 8.2/buster/nginx/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv sh

sed -Ei 's/^(access.log)/; \1/' /usr/local/etc/php-fpm.d/*.conf

exec php-fpm
